# deployment/init.d/legal-cms-md-backend
#!/bin/sh
### BEGIN INIT INFO
# Provides:          legal-cms-md-backend
# Required-Start:    $local_fs $network $postgresql
# Required-Stop:     $local_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Legal CMS-MD Backend Service
# Description:       FastAPI backend для Legal CMS-MD
### END INIT INFO

# Путь к проекту
PROJECT_DIR="/home/maimik/Projects/Legal_CMS-MD"
VENV_PATH="$PROJECT_DIR/backend/venv"
APP_PATH="$PROJECT_DIR/backend"
PID_FILE="/var/run/legal-cms-md-backend.pid"
LOG_FILE="/var/log/legal-cms-md-backend.log"
USER="maimik"

# Команда запуска
START_CMD="$VENV_PATH/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000"

case "$1" in
    start)
        echo "Starting Legal CMS-MD Backend..."

        # Проверить что процесс не запущен
        if [ -f "$PID_FILE" ]; then
            PID=$(cat "$PID_FILE")
            if kill -0 "$PID" 2>/dev/null; then
                echo "Backend already running (PID: $PID)"
                exit 1
            else
                rm -f "$PID_FILE"
            fi
        fi

        # Запустить backend
        cd "$APP_PATH"
        su - "$USER" -c "cd '$APP_PATH' && $START_CMD >> $LOG_FILE 2>&1 &"

        # Сохранить PID
        sleep 2
        PID=$(pgrep -f "uvicorn app.main:app")
        if [ -n "$PID" ]; then
            echo "$PID" > "$PID_FILE"
            echo "Backend started (PID: $PID)"
        else
            echo "Failed to start backend"
            exit 1
        fi
        ;;

    stop)
        echo "Stopping Legal CMS-MD Backend..."

        if [ -f "$PID_FILE" ]; then
            PID=$(cat "$PID_FILE")
            if kill -0 "$PID" 2>/dev/null; then
                kill "$PID"
                rm -f "$PID_FILE"
                echo "Backend stopped"
            else
                echo "Backend not running"
                rm -f "$PID_FILE"
            fi
        else
            echo "PID file not found"
        fi
        ;;

    restart)
        $0 stop
        sleep 2
        $0 start
        ;;

    status)
        if [ -f "$PID_FILE" ]; then
            PID=$(cat "$PID_FILE")
            if kill -0 "$PID" 2>/dev/null; then
                echo "Backend is running (PID: $PID)"
            else
                echo "Backend is not running (stale PID file)"
            fi
        else
            echo "Backend is not running"
        fi
        ;;

    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 1
        ;;
esac

exit 0
